# bin/mqtt/mqtt-config-watcher.service
[Unit]
Description=LoxBerry MQTT config watcher for general.json
After=network-online.target local-fs.target
Wants=network-online.target

[Service]
Type=simple
User=loxberry
Group=loxberry
ExecStart=/usr/bin/php REPLACELBPBINDIR/mqtt/mqtt_config_watcher.php
Restart=always
RestartSec=5
Nice=5
UMask=0007
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target