# bin/mqtt/mqtt-service-tts.service
[Unit]
Description=MQTT Subscriber (T2S)
After=network-online.target
Wants=network-online.target

# Start-Limits -> hierher (nicht in [Service])
StartLimitIntervalSec=60
StartLimitBurst=5
StartLimitAction=none

[Service]
Type=simple
User=loxberry
Group=loxberry
ExecStart=/usr/bin/php /opt/loxberry/bin/plugins/text2speech/mqtt/mqtt-subscribe.php
Restart=always
RestartSec=2
# (keine StartLimit* hier!)

[Install]
WantedBy=multi-user.target
