# bin/mqtt/mqtt-handshake.service
# /etc/systemd/system/mqtt-handshake.service

[Unit]
Description=T2S MQTT Handshake Listener
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=loxberry
Group=loxberry
WorkingDirectory=/opt/loxberry/bin/plugins/text2speech/mqtt
ExecStart=/usr/bin/php /opt/loxberry/bin/plugins/text2speech/mqtt/mqtt-handshake.php
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
