# bin/mqtt/mqtt-service-tts.service
[Unit]
Description=MQTT Subscriber
After=network.target

[Service]
ExecStart=/usr/bin/php REPLACELBPBINDIR/mqtt/mqtt-subscribe.php
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5
User=loxberry
WorkingDirectory=REPLACELBPBINDIR/mqtt
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target